# templates/gentoo/moxi-server-init.erb
#!/sbin/runscript

depend() {
    need net
    # Same as memcached, it makes sense
    after ntp-client ntpd rdate openrdate adjtimex hwclock
}

start() {
    ebegin "Starting moxi"
    if [ -f "${PIDFILE}" ]; then
        einfo "  Removing stale pidfile ${PIDFILE}"
        kill -9 `cat ${PIDFILE}` &>/dev/null
        rm -f "${PIDFILE}"
    fi
    /sbin/start-stop-daemon --start --pidfile "${PIDFILE}" \
        --exec "${BASEDIR}/bin/moxi" \
        --user "moxi" -- \
        -d -P "${PIDFILE}" ${OPTIONS} \
        -Z "${BASEDIR}/etc/moxi.cfg" \
        -z "${BASEDIR}/etc/moxi-cluster.cfg"
    eend $?
}

stop() {
    ebegin "Stopping moxi"
    start-stop-daemon --stop --quiet --pidfile "${PIDFILE}" \
        --exec "${BASEDIR}/bin/moxi"
    [ -f "${PIDFILE}" ] && kill -9 `cat ${PIDFILE}` &>/dev/null
    rm -f "${PIDFILE}"
    eend $?
}

